# Firmware/usr/lib/opkg/info/lvgl_klipper.list
/etc/sysConfig/UIResource/K1/img_refresh.png
/etc/sysConfig/defData/error_code_map.json
/etc/sysConfig/defData/material/material_database.json
/etc/sysConfig/UIResource/K1/img_new_chamber.png
/usr/bin/upgrade-server
/etc/sysConfig/UIResource/K1/lock.png
/etc/sysConfig/UIResource/K1/state_normal.png
/etc/sysConfig/UIResource/K1_Max/img_humidity.png
/etc/sysConfig/UIResource/K1/zoffset_Add_Normal.png
/etc/sysConfig/UIResource/K1/selfchecking_1.png
/etc/sysConfig/UIResource/K1/img_test_tip_k2_machine.png
/etc/sysConfig/UIResource/K1/fan_schematic_k2pro_longword.png
/etc/sysConfig/UIResource/K1/multicolor_box_k2pro.png
/etc/sysConfig/UIResource/K1/img_f028_fan_schematic.png
/etc/sysConfig/UIResource/K1/belt_install_left_tool.png
/etc/sysConfig/UIResource/K1/img_cfs_humidity_normal.png
/etc/sysConfig/audio/chinese/1.button_beep.wav
/etc/sysConfig/UIResource/K1/zoffset_Sub_Normal.png
/etc/sysConfig/UIResource/K1/img_new_resume.png
/etc/sysConfig/UIResource/K1/img_new_account_mask.png
/etc/sysConfig/UIResource/K1/NotVisible.png
/etc/sysConfig/UIResource/K1_Max/K2_Checktips.png
/etc/sysConfig/UIResource/K1/SortDownUnSelect.png
/etc/sysConfig/UIResource/K1/excluceObjectIcon.png
/etc/sysConfig/UIResource/K1/errorTips.png
/etc/sysConfig/UIResource/K1/img_new_exclude_objects.png
/etc/sysConfig/defData/file_preview.png
/etc/sysConfig/UIResource/K1/switch_axis_test_option_dir.png
/etc/sysConfig/defData/folder_icon.png
/etc/sysConfig/UIResource/K1/SortDownSelect.png
/etc/sysConfig/UIResource/K1/img_fila_k2pro.png
/usr/bin/audio-server
/etc/sysConfig/UIResource/K1/localexporttips.png
/etc/sysConfig/UIResource/K1/img_humidity_arrow.png
/usr/bin/Monitor
/etc/sysConfig/UIResource/K1/belt_uninstall_right_tool.png
/etc/sysConfig/UIResource/K1/switch_axis_test_step_dir.png
/etc/init.d/app
/etc/sysConfig/UIResource/K1/gcode_weight_dir_flag.png
/etc/sysConfig/UIResource/K1/img_boot0.png
/etc/sysConfig/audio/chinese/1.sound_button_smart.wav
/etc/sysConfig/defData/material/F008_material_database.json
/etc/sysConfig/defData/cert/server.crt
/etc/sysConfig/UIResource/K1/img_box_guide.png
/etc/sysConfig/UIResource/K1/selfchecking_2.png
/etc/sysConfig/defData/Auto_pressure_advance_testpadvance_product_test.gcode
/etc/sysConfig/UIResource/K1/bed.png
/usr/bin/web-server
/etc/sysConfig/UIResource/K1_Max/img_setup.png
/etc/sysConfig/UIResource/K1/img_spool_holder_k2pro.png
/etc/sysConfig/UIResource/K1/img_warn.png
/etc/sysConfig/UIResource/K1/fan_schematic_k2_shortword.png
/etc/sysConfig/UIResource/K1/time_used_1.png
/etc/sysConfig/defData/laser_offset_correction_0.2mm.gcode
/etc/sysConfig/UIResource/K1/cut_test_vector.png
/etc/sysConfig/UIResource/K1/usedtimeicon.png
/etc/sysConfig/UIResource/K1/fan_schematic_k2_s1_longword.png
/etc/sysConfig/defData/cert/server.key
/etc/sysConfig/UIResource/K1/img_fila_k2_s1.png
/etc/sysConfig/UIResource/K1/TilteLeftstyle.png
/etc/sysConfig/UIResource/K1/Subtract.png
/etc/sysConfig/UIResource/K1/state_error.png
/etc/sysConfig/UIResource/K1/TitleRightstyle.png
/etc/sysConfig/UIResource/K1/objectStateFocus.png
/etc/sysConfig/UIResource/K1/selfchecking_3.png
/etc/sysConfig/UIResource/K1/img_qrcode_tips.png
/etc/sysConfig/UIResource/K1_Max/K1MAX_Machine.png
/etc/sysConfig/UIResource/K1_Max/img_nozzle.png
/etc/sysConfig/UIResource/K1_Max/img_eye.png
/etc/sysConfig/UIResource/K1/print_success.png
/etc/sysConfig/UIResource/K1/ui_img_chamber_heating.png
/etc/sysConfig/UIResource/K1/Visible.png
/etc/sysConfig/UIResource/K1/img_k2plus_machine.png
/etc/sysConfig/UIResource/K1_Max/img_vector_1120.png
/etc/sysConfig/UIResource/K1/K1_CheckTips.png
/etc/sysConfig/UIResource/K1/fila_used_1.png
/etc/sysConfig/UIResource/K1/EstateOut.png
/etc/sysConfig/UIResource/K1/belt_dir.png
/etc/sysConfig/UIResource/K1/img_spool_holder_k2_s1.png
/etc/sysConfig/UIResource/K1/img_down_vector.png
/etc/sysConfig/UIResource/K1/img_new_binding_user_mask.png
/etc/sysConfig/UIResource/K1/img_cfs_sn_vector.png
/etc/sysConfig/defData/file_print_photo.png
/etc/sysConfig/defData/line_width_height.gcode
/etc/sysConfig/UIResource/K1/zoffset_Add_Disable.png
/etc/sysConfig/UIResource/K1/gcode_weight_flag.png
/etc/sysConfig/UIResource/K1/layers.png
/etc/sysConfig/audio/chinese/1.sound_button_crisp.wav
/etc/sysConfig/UIResource/K1/img_new_pause.png
/usr/bin/burn-server
/etc/sysConfig/UIResource/K1/ui_img_dir_down.png
/etc/sysConfig/UIResource/K1/SortIcon.png
/etc/sysConfig/UIResource/K1/img_boot0_K2.png
/etc/sysConfig/UIResource/K1/zoffset_Sub_Disable.png
/etc/sysConfig/UIResource/K1/objectStateCancle.png
/etc/sysConfig/UIResource/K1/objectStatePrinting.png
/usr/share/ca.crt
/etc/sysConfig/UIResource/K1/nuzzle.png
/etc/sysConfig/UIResource/K1/multicolor_box_k2.png
/etc/sysConfig/UIResource/K1_Max/img_polygon_25.png
/etc/sysConfig/UIResource/K1/img_new_wifi_high.png
/etc/sysConfig/UIResource/K1/color_fila_next.png
/etc/sysConfig/UIResource/K1/img_refill_dir.png
/etc/sysConfig/UIResource/K1/img_k2_s1_boot_machine.png
/etc/sysConfig/UIResource/K1/belt_uninstall_left_tool.png
/etc/sysConfig/UIResource/K1/img_humidity_normal_rangle.png
/etc/sysConfig/UIResource/K1/img_new_wifi_mid.png
/etc/sysConfig/UIResource/K1/cameraicon.png
/etc/sysConfig/UIResource/K1/EstateIn.png
/etc/sysConfig/UIResource/K1_Max/img_k2_s1_check_tips.png
/etc/sysConfig/UIResource/K1/img_cfs_box.png
/etc/sysConfig/defData/file_icon.png
/usr/bin/wifi-server
/etc/sysConfig/defData/material/F025_material_database.json
/etc/sysConfig/UIResource/K1/lightclose.png
/etc/sysConfig/UIResource/K1/img_k2_s1_machine.png
/etc/sysConfig/UIResource/K1/img_stop_fan.png
/etc/sysConfig/UIResource/K1/icon_chamber_heating.png
/etc/sysConfig/UIResource/K1/icon_chamber.png
/etc/sysConfig/UIResource/K1/uploadtips.png
/etc/sysConfig/UIResource/K1/udiskicon.png
/etc/sysConfig/UIResource/K1/img_test_tip_k2_s1_mchine.png
/etc/sysConfig/UIResource/K1/wifi_state_1.png
/etc/sysConfig/defData/Auto_Em.gcode
/etc/sysConfig/UIResource/K1/img_new_stop.png
/etc/sysConfig/UIResource/K1/img_new_bed.png
/etc/sysConfig/UIResource/K1/silentmode.png
/etc/sysConfig/audio/chinese/1.sound_button_technology.wav
/etc/sysConfig/defData/material/F012_material_database.json
/etc/sysConfig/UIResource/K1/selected_state.png
/etc/sysConfig/defData/laser_offset_correction.gcode
/etc/sysConfig/UIResource/K1/EstateRunoutIn.png
/etc/sysConfig/UIResource/K1/lightopen.png
/etc/sysConfig/UIResource/K1/EstateNormal.png
/etc/sysConfig/UIResource/K1/color_fila_prev.png
/etc/sysConfig/UIResource/K1/warningTips.png
/etc/sysConfig/defData/cert/ca.crt
/etc/sysConfig/UIResource/K1/axiscoolpagefan_open.png
/etc/sysConfig/UIResource/K1/img_new_wifi_low.png
/etc/sysConfig/defData/material/F021_material_database.json
/etc/sysConfig/UIResource/K1/checking.gif
/etc/sysConfig/UIResource/K1/K1_Machine.png
/etc/sysConfig/UIResource/K1/SwitchAxisStep1.png
/etc/sysConfig/UIResource/K1/belt_install_right_tool.png
/etc/sysConfig/UIResource/K1/img_boot4.png
/etc/sysConfig/UIResource/K1/img_f028_spool_holder.png
/etc/sysConfig/UIResource/K1/nozzle_fila_empty.png
/etc/sysConfig/UIResource/K1/auto_refill_prev_page.png
/etc/sysConfig/UIResource/K1/img_fila_k2plus.png
/etc/sysConfig/UIResource/K1/img_test_tip_f028_machine.png
/etc/sysConfig/UIResource/K1/img_new_light.png
/etc/sysConfig/UIResource/K1/wifi_state_2.png
/etc/sysConfig/UIResource/K1/homeZ.png
/etc/sysConfig/audio/chinese/1.sound_button_short.wav
/etc/sysConfig/UIResource/K1/SortUpUnSelect.png
/etc/sysConfig/UIResource/K1/speed_mode.png
/etc/sysConfig/UIResource/K1/img_polygon_border.png
/etc/sysConfig/UIResource/K1/img_spool_holder_k2.png
/etc/sysConfig/UIResource/K1/img_multi_color_guide_3.png
/usr/bin/display-server
/etc/sysConfig/UIResource/K1/filamentTip.png
/etc/sysConfig/UIResource/K1/fan_schematic.png
/etc/sysConfig/UIResource/K1/img_humidity_little_height_rangle.png
/etc/sysConfig/UIResource/K1/SwitchAxisStep2.png
/etc/sysConfig/UIResource/K1/img_in_progress.gif
/etc/sysConfig/UIResource/K1/img_humidity_height_rangle.png
/etc/sysConfig/UIResource/K1/img_test_tip_k2plus_machine.png
/etc/sysConfig/UIResource/K1/img_logo_animation.gif
/etc/sysConfig/UIResource/K1_Max/K2Plus_Checktips.png
/etc/sysConfig/UIResource/K1_Max/K1Max_Checktips.png
/etc/sysConfig/UIResource/K1/turn_right.png
/etc/sysConfig/UIResource/K1/img_new_camera.png
/etc/sysConfig/UIResource/K1/K1C_Machine.png
/etc/sysConfig/UIResource/K1/turn_left.png
/usr/bin/master-server
/etc/sysConfig/defData/SourceHanSansCN-Normal.otf
/etc/sysConfig/UIResource/K1/wifi_state_3.png
/etc/sysConfig/defData/Auto_pressure_advance_multi.gcode
/etc/sysConfig/UIResource/K1_Max/img_edit.png
/etc/sysConfig/UIResource/K1/SortSelect.png
/etc/sysConfig/UIResource/K1/Add.png
/etc/sysConfig/UIResource/K1/img_cfs_hub_indicator_light.png
/etc/sysConfig/UIResource/K1/hotendicon.png
/etc/sysConfig/UIResource/K1/auto_refill_next_page.png
/usr/bin/app-server
/etc/sysConfig/UIResource/K1/img_fila_f028.png
/etc/sysConfig/UIResource/K1/belt_init.png
/etc/sysConfig/UIResource/K1/img_rfid_flag.png
/etc/sysConfig/UIResource/K1/ui_img_dir_up.png
/etc/sysConfig/UIResource/K1_Max/img_f028_Checktips.png
/etc/sysConfig/UIResource/K1/img_spool_holder.png
/etc/sysConfig/UIResource/K1/EstateRunoutOut.png
/etc/sysConfig/UIResource/K1/img_cfs_hub_light.png
/etc/sysConfig/UIResource/K1/nozzle_fila_full.png
/etc/sysConfig/UIResource/K1_Max/K2Pro_Checktips.png
/usr/share/server.crt
/etc/sysConfig/UIResource/K1/img_cfs_humidity_height.png
/etc/sysConfig/UIResource/K1/img_k2pro_machine.png
/etc/sysConfig/UIResource/K1/img_f028_boot_machine.png
/etc/sysConfig/UIResource/K1/zoffset_Add_Focus.png
/etc/sysConfig/UIResource/K1/print_completed.png
/etc/sysConfig/UIResource/K1/zoffset_Sub_Focus.png
/etc/sysConfig/UIResource/K1_Max/img_polygon.png
/etc/sysConfig/UIResource/K1/fan_schematic_k2_longword.png
/etc/sysConfig/UIResource/K1/EstateRunout.png
/etc/sysConfig/UIResource/K1/img_small_default_user_avatar.png
/etc/sysConfig/UIResource/K1/wait.png
/etc/sysConfig/UIResource/K1/SwitchAxisStep3.png
/etc/sysConfig/UIResource/K1/img_guide_page1.png
/etc/sysConfig/defData/Auto_pressure_advance_testpadvance.gcode
/etc/sysConfig/UIResource/K1/stop_print_warning.png
/etc/sysConfig/UIResource/K1/img_cfs_humidity_little_height.png
/etc/sysConfig/UIResource/K1/SortUpSelect.png
/etc/sysConfig/UIResource/K1/img_home_user_mask.png
/etc/sysConfig/UIResource/K1/img_f028_machine.png
/etc/sysConfig/UIResource/K1/remainingtimeicon.png
/etc/sysConfig/defData/laser_test_line.gcode
/etc/sysConfig/UIResource/K1/selfchecking_0.png
/etc/sysConfig/UIResource/K1/fan_schematic_k2pro_shortword.png
/etc/sysConfig/UIResource/K1/wifi_state_4.png
/etc/sysConfig/UIResource/K1/axiscoolpagefan.png
/etc/sysConfig/UIResource/K1/img_fila_k2.png
/usr/share/server.key
/etc/sysConfig/UIResource/K1/objectStateNormal.png
/etc/sysConfig/UIResource/K1/img_test_tip_k2pro_machine.png
/etc/sysConfig/UIResource/K1/img_k2_machine.png
/etc/sysConfig/UIResource/K1/bedicon.png
/etc/sysConfig/UIResource/K1/gcode_default_small.png
/etc/sysConfig/defData/auto_laser_test.gcode
/etc/sysConfig/UIResource/K1/img_new_hotend.png
